# storage/usr/share/ati/drv.list
6f10a2bf34fe188ae357ed3b537b4fa4  /usr/share/doc/ati/user-manual/index.html
1d6f434a8e99b435e7ddb49e2e5fcdf7  /usr/share/doc/ati/examples/etc/acpi/ati-powermode.sh
fe6c2398b1890a0b0ecbbf2617471827  /usr/share/doc/ati/examples/etc/acpi/events/a-ac-aticonfig
fdc02818706c3ae5062a9add9125448d  /usr/share/doc/ati/examples/etc/acpi/events/a-lid-aticonfig
9cb2a07899a47e324120a2da39615e60  /usr/share/doc/ati/examples/etc/init.d/atieventsd.sh
4c738c120bb7447bbf33c5c4de2f56b7  /usr/share/doc/ati/articles/1gbhang.html
2518977c826236be43a416190179197d  /usr/share/doc/ati/articles/4461.html
4f6468d5e8b8dbc7a08dc8c25517941f  /usr/share/doc/ati/articles/4462.html
a0dc065bd868ca072100a181045ab0d1  /usr/share/doc/ati/articles/4463.html
758e6ce456114da4ff6fa5edd7baefa9  /usr/share/doc/ati/articles/4464.html
71e699def415041aca04dea8f4484e82  /usr/share/doc/ati/articles/4469.html
b3db8a0a7498dea2af06e999575a2d97  /usr/share/doc/ati/articles/4470.html
a04f81945b48cb5be3f291fa95796b05  /usr/share/doc/ati/articles/4475.html
cb7b976a5a36d688af8b025713478262  /usr/share/doc/ati/articles/4478.html
c621e5eccdedb251fcab1b2e20d6d5d2  /usr/share/doc/ati/articles/4479.html
c13b19e7c267391291ae2bc18d196c88  /usr/share/doc/ati/articles/4480.html
67d83914c446e6ec7a69564939794b34  /usr/share/doc/ati/articles/4481.html
09fd9ff03df6130efad58b8906b03eca  /usr/share/doc/ati/articles/4482.html
a1c14596671271b730b1e4a63b6a3c65  /usr/share/doc/ati/articles/4483.html
9801b27dc9371fa024b7e8c02219a848  /usr/share/doc/ati/articles/4484.html
c3974b76fb457adbb6e44cf2d8119628  /usr/share/doc/ati/articles/4485.html
77299d51e5572f82b5c1d2041e7bf228  /usr/share/doc/ati/articles/corruptstereo.html
86ea2b26ae456d7586fb44ee6d3f24f2  /usr/share/doc/ati/articles/corruptvtswitch.html
4e4ce6097704fa6eb36bc97e462da296  /usr/share/doc/ati/articles/devshm.html
33c6e2ae36cb4b9cfab15a158be28eb2  /usr/share/doc/ati/articles/dga3dhang.html
caaf3ffa4a9cf527d24109d45d20ebfb  /usr/share/doc/ati/articles/doom3corrupt.html
376c128b48ba09657b33b97f108a29e1  /usr/share/doc/ati/articles/dualheadvideo.html
d5bd2ffc7b0c3c93474246943af68962  /usr/share/doc/ati/articles/laptopsuspend.html
bbdf0c411b334db6ddab6a4ea511f6c2  /usr/share/doc/ati/articles/missingdrmheaders.html
a6c1f139a03e3bcacff3e3e7a1cbc25d  /usr/share/doc/ati/articles/mousecursorhang.html
b334fa804ab0fcff0e33cb59d800c5b3  /usr/share/doc/ati/articles/no3d-aiw8500dv.html
61789a8de4dd301de4c3ffa714e15199  /usr/share/doc/ati/articles/no3d-kt400.html
166c7e32b05aa78a09f12261b294ade9  /usr/share/doc/ati/articles/nomembercount.html
871d428f2eb939855beb80b289f1c8ed  /usr/share/doc/ati/articles/pcie3dmemoryleak.html
158cf676d4b2b83379dd2c0c4910b46d  /usr/share/doc/ati/articles/r420blankdisplay.html
0bc580f23c5d4380e71e265dcc5182eb  /usr/share/doc/ati/articles/rv280dviblankdisplay.html
1c5005ecc2f4d3f733b54c49369ad028  /usr/share/doc/ati/articles/rv350springdale.html
9e8651a0f6aef4cc48110b558c834e90  /usr/share/doc/ati/articles/secondheadcorruption.html
9ffa905d815ccc63dcc2ecd94c6c840b  /usr/share/doc/ati/articles/xf86_enodev.html
20662519b73112451bd743aac160de80  /usr/share/doc/ati/articles/xrestartpcie.html
8cc6585ddad4a8a81e1067f0787a8d4a  /usr/share/doc/ati/articles/xvsatshift.html
f917d21a49d86498b8a16c9f2391e8e7  /usr/share/doc/ati/configure.html
552fbf80507d31057a97db25ffc684af  /usr/share/doc/ati/driverfaq.html
1d1a864752156aad073908e671d3c896  /usr/share/doc/ati/index.html
1676e09231e04265ce16e5d2e00779f0  /usr/share/doc/ati/installer.html
f27006cf380c2002770188bfaa92a08a  /usr/share/doc/ati/issues.html
9e470e3116155a0e2cb36f962803e56c  /usr/share/doc/ati/linuxfaq.html
5f834bc04e226b9cdb9004ff4b182c43  /usr/share/doc/ati/tips-linux.html
b3868935bd730a291057da00a749b537  /usr/share/man/man8/atieventsd.8
87e13f477dd307ac65d8eed4db2c1020  /usr/lib64/xorg/modules/linux/libfglrxdrm.so
91cb8979d7ffbd0bd67c29ef490d59dc  /usr/lib64/xorg/modules/drivers/fglrx_drv.so
02a76b2bc5d317176d70b9574a4d62ef  /usr/lib64/xorg/modules/glesx.so
d216723c5c69d89a5c972775bd774204  /usr/lib64/xorg/modules/amdxmm.so
d1d504bae5763b383d9d0947447974d0  /usr/lib/libAMDXvBA.cap
fab7a9b4656a4a9f369ab211ad5d3295  /usr/lib/libAMDXvBA.so.1.0
5591700bb2fc31ee2e61203f2d23c29b  /usr/lib/libatiadlxx.so
1cea1dd0e449e4c07daff54327904107  /usr/lib/libaticalcl.so
034287d832ef2e96d692914f6469c9db  /usr/lib/libaticaldd.so
34757979118625b9f0fde391a8738049  /usr/lib/libaticalrt.so
08d9aafe409c74a0dc91824dc0ce2bdb  /usr/lib/libatiuki.so.1.0
9fa2d61df1604a5aac96d44f24771d51  /usr/lib/libfglrx_dm.a
3e4aaafea3141fb48c6c7591874f81c3  /usr/lib/libfglrx_dm.so.1.0
c06ed7f4f874b0ddc86485ed93425880  /usr/lib/libXvBAW.so.1.0
c31338fb7ee3c7c8534499fdfac04bd0  /usr/lib/fglrx/fglrx-libGL.so.1.2
19e0f403e86a841764a98987d9fed81d  /usr/lib/dri/fglrx_dri.so
1c473e56bfb5e1c015344931fc855580  /usr/lib64/xorg/modules/extensions/fglrx/fglrx-libglx.so
d1d504bae5763b383d9d0947447974d0  /usr/lib64/libAMDXvBA.cap
6ea881f7095b3c80595004c6eee7eeb6  /usr/lib64/libAMDXvBA.so.1.0
19f56b8289c0e6f3f65166f67810a98e  /usr/lib64/libatiadlxx.so
28b8990cd6f7b4aab57efc85f9bd58da  /usr/lib64/libaticalcl.so
e5371ad3d7a0809db02ab59eaa7056c7  /usr/lib64/libaticaldd.so
d02710cc88f82e508b4e1adce2d753cb  /usr/lib64/libaticalrt.so
d972f101846af7b38d215817d3bc6f42  /usr/lib64/libatiuki.so.1.0
d656d074ba4009a56287462ceeddb032  /usr/lib64/libfglrx_dm.a
070ca71078f0040cea382a376bfc5d3a  /usr/lib64/libfglrx_dm.so.1.0
4b0e98854b77995ec9fd13c13fd18fe7  /usr/lib64/libXvBAW.so.1.0
0e9fd011f079fcb07994d19fcbc8fd3c  /usr/lib64/fglrx/fglrx-libGL.so.1.2
0be401cde0f355d7043dd391bb76e7ff  /usr/lib64/dri/fglrx_dri.so
9cf308d2885d31c991ce4e6a28e9cb87  /usr/include/GL/glATI.h
c79c61e102fe88ed0c23c1674d1f5427  /usr/include/GL/glxATI.h
0d1de2427998218fd2130002a903f547  /usr/include/ATI/GL/glxext.h
8f29336cecf9aa2a57b3c87436fe15b8  /usr/include/ATI/GL/glx.h
036161bb2511b79a8860b00cf836a690  /usr/bin/fgl_glxgears
54f179e68fac597e2f65cfe7d9cf6685  /usr/bin/fglrxinfo
19bd3024b5f5ae1b6d5d6ac97c052def  /usr/bin/aticonfig
687094f13a99e39bb923bb9093c6cbb9  /usr/bin/atiodcli
0c3e85f6a66a5f01899502bbc112d667  /usr/bin/atiode
f1ee8fc0e0b3e485eb67149948070c0d  /usr/src/ati/fglrx_sample_source.tgz
3973e1bbc221ee5248dbee03a5362f31  /usr/sbin/amdnotifyui
24206bfadcdd74084ed66eca3be845e6  /usr/sbin/atieventsd
b17b468a89e9188b779aadb8efa7c501  /usr/sbin/atigetsysteminfo.sh
ba1e80731a4c60ae929e0201abc6a5bc  /etc/ati/amdpcsdb.default
1ece1c6faeb69dcf1cb1ae77dd141b5a  /etc/ati/atiogl.xml
986a6d490e6fcb5faf915aec9857189f  /etc/ati/authatieventsd.sh
4a3697db897594fc136fef9359c93d16  /etc/ati/control
16761230b54565b803b8ff6021a9d071  /etc/ati/logo_mask.xbm.example
5c0938e882815a29cbd7e9a853957575  /etc/ati/logo.xbm.example
d3036ceec6741118dca940369944c3e7  /etc/ati/signature
19e0f403e86a841764a98987d9fed81d  /usr/X11R6/lib/modules/dri/fglrx_dri.so
0be401cde0f355d7043dd391bb76e7ff  /usr/X11R6/lib64/modules/dri/fglrx_dri.so
6f10a2bf34fe188ae357ed3b537b4fa4  /usr/share/doc/ati/user-manual/index.html
9cb2a07899a47e324120a2da39615e60  /usr/share/doc/ati/examples/etc/init.d/atieventsd.sh
1d6f434a8e99b435e7ddb49e2e5fcdf7  /usr/share/doc/ati/examples/etc/acpi/ati-powermode.sh
fdc02818706c3ae5062a9add9125448d  /usr/share/doc/ati/examples/etc/acpi/events/a-lid-aticonfig
fe6c2398b1890a0b0ecbbf2617471827  /usr/share/doc/ati/examples/etc/acpi/events/a-ac-aticonfig
4c738c120bb7447bbf33c5c4de2f56b7  /usr/share/doc/ati/articles/1gbhang.html
2518977c826236be43a416190179197d  /usr/share/doc/ati/articles/4461.html
4f6468d5e8b8dbc7a08dc8c25517941f  /usr/share/doc/ati/articles/4462.html
a0dc065bd868ca072100a181045ab0d1  /usr/share/doc/ati/articles/4463.html
758e6ce456114da4ff6fa5edd7baefa9  /usr/share/doc/ati/articles/4464.html
71e699def415041aca04dea8f4484e82  /usr/share/doc/ati/articles/4469.html
b3db8a0a7498dea2af06e999575a2d97  /usr/share/doc/ati/articles/4470.html
a04f81945b48cb5be3f291fa95796b05  /usr/share/doc/ati/articles/4475.html
cb7b976a5a36d688af8b025713478262  /usr/share/doc/ati/articles/4478.html
c621e5eccdedb251fcab1b2e20d6d5d2  /usr/share/doc/ati/articles/4479.html
c13b19e7c267391291ae2bc18d196c88  /usr/share/doc/ati/articles/4480.html
67d83914c446e6ec7a69564939794b34  /usr/share/doc/ati/articles/4481.html
09fd9ff03df6130efad58b8906b03eca  /usr/share/doc/ati/articles/4482.html
a1c14596671271b730b1e4a63b6a3c65  /usr/share/doc/ati/articles/4483.html
9801b27dc9371fa024b7e8c02219a848  /usr/share/doc/ati/articles/4484.html
c3974b76fb457adbb6e44cf2d8119628  /usr/share/doc/ati/articles/4485.html
77299d51e5572f82b5c1d2041e7bf228  /usr/share/doc/ati/articles/corruptstereo.html
86ea2b26ae456d7586fb44ee6d3f24f2  /usr/share/doc/ati/articles/corruptvtswitch.html
4e4ce6097704fa6eb36bc97e462da296  /usr/share/doc/ati/articles/devshm.html
33c6e2ae36cb4b9cfab15a158be28eb2  /usr/share/doc/ati/articles/dga3dhang.html
caaf3ffa4a9cf527d24109d45d20ebfb  /usr/share/doc/ati/articles/doom3corrupt.html
376c128b48ba09657b33b97f108a29e1  /usr/share/doc/ati/articles/dualheadvideo.html
d5bd2ffc7b0c3c93474246943af68962  /usr/share/doc/ati/articles/laptopsuspend.html
bbdf0c411b334db6ddab6a4ea511f6c2  /usr/share/doc/ati/articles/missingdrmheaders.html
a6c1f139a03e3bcacff3e3e7a1cbc25d  /usr/share/doc/ati/articles/mousecursorhang.html
b334fa804ab0fcff0e33cb59d800c5b3  /usr/share/doc/ati/articles/no3d-aiw8500dv.html
61789a8de4dd301de4c3ffa714e15199  /usr/share/doc/ati/articles/no3d-kt400.html
166c7e32b05aa78a09f12261b294ade9  /usr/share/doc/ati/articles/nomembercount.html
871d428f2eb939855beb80b289f1c8ed  /usr/share/doc/ati/articles/pcie3dmemoryleak.html
158cf676d4b2b83379dd2c0c4910b46d  /usr/share/doc/ati/articles/r420blankdisplay.html
0bc580f23c5d4380e71e265dcc5182eb  /usr/share/doc/ati/articles/rv280dviblankdisplay.html
1c5005ecc2f4d3f733b54c49369ad028  /usr/share/doc/ati/articles/rv350springdale.html
9e8651a0f6aef4cc48110b558c834e90  /usr/share/doc/ati/articles/secondheadcorruption.html
9ffa905d815ccc63dcc2ecd94c6c840b  /usr/share/doc/ati/articles/xf86_enodev.html
20662519b73112451bd743aac160de80  /usr/share/doc/ati/articles/xrestartpcie.html
8cc6585ddad4a8a81e1067f0787a8d4a  /usr/share/doc/ati/articles/xvsatshift.html
f917d21a49d86498b8a16c9f2391e8e7  /usr/share/doc/ati/configure.html
552fbf80507d31057a97db25ffc684af  /usr/share/doc/ati/driverfaq.html
1d1a864752156aad073908e671d3c896  /usr/share/doc/ati/index.html
1676e09231e04265ce16e5d2e00779f0  /usr/share/doc/ati/installer.html
f27006cf380c2002770188bfaa92a08a  /usr/share/doc/ati/issues.html
9e470e3116155a0e2cb36f962803e56c  /usr/share/doc/ati/linuxfaq.html
5f834bc04e226b9cdb9004ff4b182c43  /usr/share/doc/ati/tips-linux.html
b3868935bd730a291057da00a749b537  /usr/share/man/man8/atieventsd.8
86bfeafa487cbbce9c5e1491a39965f1  /usr/lib64/xorg/modules/linux/libfglrxdrm.so
19dc8e1285038ce25809a062528dc703  /usr/lib64/xorg/modules/drivers/fglrx_drv.so
68cb52e10b1e2b0e3019aa9be3ae6f75  /usr/lib64/xorg/modules/glesx.so
9f1a4d8c5d31668bbc8e66d985555fb9  /usr/lib64/xorg/modules/amdxmm.so
c9e845fc213b4dee4b2dacd412390b44  /usr/lib/libAMDXvBA.cap
6264b9acaa9781901f2da0cbc283f90c  /usr/lib/libAMDXvBA.so.1.0
f4ecffea3dc0af198b664710dd628f04  /usr/lib/libatiadlxx.so
1cea1dd0e449e4c07daff54327904107  /usr/lib/libaticalcl.so
79767e4e35c3422adc6878949469b8da  /usr/lib/libaticaldd.so
34757979118625b9f0fde391a8738049  /usr/lib/libaticalrt.so
08d9aafe409c74a0dc91824dc0ce2bdb  /usr/lib/libatiuki.so.1.0
1b0e571295c6691c39bc14a50fac1148  /usr/lib/libfglrx_dm.a
3e4aaafea3141fb48c6c7591874f81c3  /usr/lib/libfglrx_dm.so.1.0
c06ed7f4f874b0ddc86485ed93425880  /usr/lib/libXvBAW.so.1.0
d42540a1d70e8250aff60accaab06806  /usr/lib/fglrx/libGL.so.1.2
af196c98a29a5c54e9b995db0e2840bd  /usr/lib/dri/fglrx_dri.so
63703a474a50508929e0ace1e565e9c2  /usr/lib64/xorg/modules/extensions/fglrx/libglx.so
c9e845fc213b4dee4b2dacd412390b44  /usr/lib64/libAMDXvBA.cap
bb87235fc58cf718c28d136f17141dc8  /usr/lib64/libAMDXvBA.so.1.0
45226810030be3d9207a4f9a07c4cd05  /usr/lib64/libatiadlxx.so
28b8990cd6f7b4aab57efc85f9bd58da  /usr/lib64/libaticalcl.so
28fc9183d49d97312b6004af90e0afb5  /usr/lib64/libaticaldd.so
d02710cc88f82e508b4e1adce2d753cb  /usr/lib64/libaticalrt.so
d972f101846af7b38d215817d3bc6f42  /usr/lib64/libatiuki.so.1.0
6fd6407707d078b80f2a21f3e3001a9a  /usr/lib64/libfglrx_dm.a
070ca71078f0040cea382a376bfc5d3a  /usr/lib64/libfglrx_dm.so.1.0
4b0e98854b77995ec9fd13c13fd18fe7  /usr/lib64/libXvBAW.so.1.0
923cea1f4969f3a2a4296c24c66dbe30  /usr/lib64/fglrx/libGL.so.1.2
b71725ad30d10b4fcb45739209db7356  /usr/lib64/dri/fglrx_dri.so
9cf308d2885d31c991ce4e6a28e9cb87  /usr/include/GL/glATI.h
c79c61e102fe88ed0c23c1674d1f5427  /usr/include/GL/glxATI.h
0d1de2427998218fd2130002a903f547  /usr/include/ATI/GL/glxext.h
8f29336cecf9aa2a57b3c87436fe15b8  /usr/include/ATI/GL/glx.h
036161bb2511b79a8860b00cf836a690  /usr/bin/fgl_glxgears
54f179e68fac597e2f65cfe7d9cf6685  /usr/bin/fglrxinfo
5808a20919da27a39e851d7c49aa72ca  /usr/bin/aticonfig
687094f13a99e39bb923bb9093c6cbb9  /usr/bin/atiodcli
0c3e85f6a66a5f01899502bbc112d667  /usr/bin/atiode
f1ee8fc0e0b3e485eb67149948070c0d  /usr/src/ati/fglrx_sample_source.tgz
3973e1bbc221ee5248dbee03a5362f31  /usr/sbin/amdnotifyui
7eb73873c4c8cd3d995d89b46560625c  /usr/sbin/atieventsd
b17b468a89e9188b779aadb8efa7c501  /usr/sbin/atigetsysteminfo.sh
860b770e98983bb8a3b7977295817ad3  /etc/ati/amdpcsdb.default
3638d76ed11785b9aa72eb9539433a63  /etc/ati/atiogl.xml
986a6d490e6fcb5faf915aec9857189f  /etc/ati/authatieventsd.sh
2d82b1de8d14d33651b5df8e92ff522a  /etc/ati/control
16761230b54565b803b8ff6021a9d071  /etc/ati/logo_mask.xbm.example
5c0938e882815a29cbd7e9a853957575  /etc/ati/logo.xbm.example
1df7525d07b4660fb3607e605e6410d1  /etc/ati/signature
af196c98a29a5c54e9b995db0e2840bd  /usr/X11R6/lib/modules/dri/fglrx_dri.so
b71725ad30d10b4fcb45739209db7356  /usr/X11R6/lib64/modules/dri/fglrx_dri.so
